# rov-kamera.service
[Unit]
Description=ROV Kamera (Picamera2 -> JPEG -> TCP Stream)
After=network.target

[Service]
Type=simple
User=ali
WorkingDirectory=/home/ali/Desktop/ROV
# -u: stdout buffering kapali, print'ler aninda journal'a gider
ExecStart=/usr/bin/python3 -u /home/ali/Desktop/ROV/kamera_yayini.py

# Cokerse otomatik tekrar basla (3sn ara, kamera serbest kalmasi icin)
Restart=always
RestartSec=3

# 5 dakikada 5'ten cok crash olursa servisi durdur (sonsuz loop'u kir)
StartLimitIntervalSec=300
StartLimitBurst=5

Environment=PYTHONUNBUFFERED=1

# Loglar journald'a yazilir
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
