# rov-beyin.service
[Unit]
Description=ROV Beyin (UDP -> Arduino Serial bridge)
# Ag + USB hazir olunca basla
After=network.target systemd-udev-settle.service
Wants=network.target

[Service]
Type=simple
User=ali
WorkingDirectory=/home/ali/Desktop/ROV

# Boot sirasinda USB enumeration sorunu olursa, Arduino'yu yeniden enumere et.
# `+` prefix: ExecStartPre root yetkisi ile calisir (User=ali gozetilmez).
# Bu, kabloyu cikar+tak ile ayni etkiyi yapar.
ExecStartPre=+/bin/bash /home/ali/Desktop/ROV/scripts/rov-usb-reset.sh

# -u: Python stdout buffering kapali, print'ler aninda journal'a gider
ExecStart=/usr/bin/python3 -u /home/ali/Desktop/ROV/rov_beyin.py

# Cokerse otomatik tekrar basla
Restart=always
RestartSec=3

# Python stdout buffering kapatma (yedek)
Environment=PYTHONUNBUFFERED=1

# Loglar journald'a yazilir (journalctl ile bakilir)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
