# rov-kamera.service
[Unit]
Description=ROV Kamera (Picamera2 -> TCP Stream)
After=network.target
Wants=network.target

[Service]
Type=simple
User=ali
WorkingDirectory=/home/ali/Desktop/ROV
# -u: stdout buffering kapali, print'ler aninda journal'a gider
ExecStart=/usr/bin/python3 -u /home/ali/Desktop/ROV/kamera_yayini.py

# Cokerse otomatik tekrar basla
Restart=always
RestartSec=3

Environment=PYTHONUNBUFFERED=1

# Loglar journald'a yazilir
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
